# attendance.service
[Unit]
Description="Attendance Bot"

[Service]
ExecStart=/usr/bin/python3 /home/pi/attendance/main.py

[Install]
WantedBy=default.target
